# config/quant-autoresearch.service
# Obsolete V1 daemon template; retained for reference only.
[Unit]
Description=Quant Autoresearch v1 daemon template (obsolete)

[Service]
Type=oneshot
ExecStart=/bin/true
Restart=no
