# config/supervisord.conf
; Obsolete V1 runtime template retained only so the old daemon surface is
; documented as unsupported in V2.
; Use uv run python cli.py backtest or follow program.md instead.

[program:quant-autoresearch]
command=/bin/sh -lc 'echo "config/supervisord.conf is obsolete for V2. Use uv run python cli.py backtest or follow program.md instead." >&2; exit 1'
directory=/Users/apple/quant-autoresearch
autostart=false
autorestart=false
startretries=0
stderr_logfile=/Users/apple/quant-autoresearch/run.err.log
stdout_logfile=/Users/apple/quant-autoresearch/run.out.log
user=apple
