# config/supervisord.conf
; Obsolete V1 supervisor template retained only for historical reference.
[program:quant-autoresearch]
command=/bin/true
directory=/tmp
autostart=false
autorestart=false
startretries=0
stderr_logfile=/dev/null
stdout_logfile=/dev/null
redirect_stderr=true
